# sddm-theme-biglinux/lib/systemd/user/sddm-wallpaper.path
[Unit]
Description=SDDM change WallPaper
ConditionUser=!root

[Path]
PathModified=%h/.config/plasma-org.kde.plasma.desktop-appletsrc
Unit=sddm-wallpaper.service

[Install]
WantedBy=default.target
